# twitter-scheduler.service
[Unit]
Description=Twitter Thread Scheduler Bot
After=network.target ollama.service

[Service]
Type=simple
User=admin_pi5
WorkingDirectory=/home/admin_pi5/twitter-scheduler
EnvironmentFile=/home/admin_pi5/twitter-scheduler/.env
ExecStart=/home/admin_pi5/twitter-scheduler/venv/bin/python bot.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
