# twitter-scheduler.service
[Unit]
Description=Twitter Thread Scheduler Bot
After=network.target

[Service]
Type=simple
User=<your-username>
WorkingDirectory=/home/<your-username>/twitter-scheduler
EnvironmentFile=/home/<your-username>/twitter-scheduler/.env
ExecStart=/home/<your-username>/twitter-scheduler/venv/bin/python bot.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
